# supervisord.conf
[supervisord]
nodaemon=true

[program:blockchain]
command=blockchain-wallet-service start -b 0.0.0.0
